# check_cmile/crontab
* * * * * /usr/src/app/main.py >> /var/log/cron.log 2>&1

